# target/linux/ramips/mt7621/base-files/etc/hotplug.d/firmware/11-mt76-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"mediatek/mt7915_eeprom_dbdc.bin")
	case "$board" in
	netis,n6)
		EEPROM=/lib/firmware/$FIRMWARE
		head -c $((0xe00)) /dev/mtd2 > $EEPROM
		printf "\x92" | \
			dd of=$EEPROM seek=$((0x190)) bs=1 conv=notrunc \
			2>/dev/null
		;;
	tplink,ex220-v2)
		EEPROM=/lib/firmware/$FIRMWARE
		mkdir -p /mnt/misc_ro
		if mount -t jffs2 -o ro /dev/mtdblock2 /mnt/misc_ro 2>/dev/null; then
		    cp /mnt/misc_ro/0x00180000 "$EEPROM" 2>/dev/null || \
		        echo "Failed to copy EEPROM file for ex220-v2"
		    umount /mnt/misc_ro
		fi
		rm -rf /mnt/misc_ro
		;;
	esac
	;;
esac
